# .local/share/systemd/user/notify-low-battery.timer
[Unit]
Description=Check for low battery

[Timer]
OnBootSec=30s
OnUnitActiveSec=90s

[Install]
WantedBy=timers.target
