# .local/share/systemd/user/notify-low-battery.service
[Unit]
Description=Notify if battery level is low

[Service]
Type=oneshot
ExecStart=%h/scripts/notify-low-battery.sh
